# config/systemd/tkd-dzwirzyno.service
[Unit]
Description=TKD Dzwirzyno Website (Deno + PocketBase)
After=network.target

[Service]
# User/Group to run the service
User=deploy
Group=deploy

# Working directory (points to the 'current' symlink)
WorkingDirectory=/var/www/tkd-dzwirzyno/current

# Environment variables
# We assume .env is symlinked into current, but Deno Fresh reads .env automatically.
# Alternatively, load them here:
EnvironmentFile=/var/www/tkd-dzwirzyno/.env

# Add deno path to environment just in case
Environment="PATH=/home/deploy/.deno/bin:/usr/local/bin:/usr/bin:/bin"

# Command to start the server
# Using 'deno task start' ensures we use the command defined in deno.json
ExecStart=/home/deploy/.deno/bin/deno task start

# Restart policy
Restart=always
RestartSec=10

# Output logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tkd-website

[Install]
WantedBy=multi-user.target
